# deployment.timer
[Unit]
Description = Timer for the service

[Timer]
OnUnitActiveSec=15sec
OnBootSec=10

[Install]
WantedBy=timers.target
