# deployment.service
[Unit]
Description= Auto website deploy

[Service]
ExecStart=/usr/bin/deploy.sh

[Install]
WantedBy=multi-user.target
